# src/data/yast-timesync.timer.erb
[Timer]
# first sync after boot
OnBootSec=1min
OnUnitActiveSec=<%= timeout %>min

[Install]
WantedBy=timers.target
